# debian/postrm
#! /bin/sh

set -e
if [ "$1" = "purge" ] ; then
        update-rc.d siguibui remove >/dev/null
fi
if [ -d /var/cache/siguibui ] ; then
	rm -Rf /var/cache/siguibui
fi