# debian/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/siguibui" ]; then
        invoke-rc.d siguibui stop || exit $?
fi
